# git-sync/failure-popup@.service
[Unit]
Description=Failure popup for %i

[Service]
Type=oneshot
ExecStart=%h/git/dotfiles/git-sync/failure.sh %i
